# build/usr/lib/systemd/system/monitor-ambient-light.service
[Unit]
Description=monitor-ambient-light

[Service]
Type=simple
ExecStart=/usr/bin/monitor-ambient-light
KillSignal=SIGINT
StandardError=null
Restart=always

[Install]
WantedBy=multi-user.target
Alias=monitor-ambient-light